# udev.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/home/jakutis/bin/udevmount"
